# build/sothothv2_agent/systemd/sothothv2_agent.service
[Unit]
Description=SothothV2 Service
Documentation=https://example.com/sothothv2-docs
After=network.target

[Service]
Type=simple
ExecStart=/sothothv2/bin/sothothv2_agent -config /sothothv2/config/sothothv2_agent.ini -foreground
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
