# build/sothothv2_agent/systemd/sothothv2_agent.service
[Unit]
Description=Sthothv2 Monitor Agent Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/monitor
ExecStart=/usr/local/bin/monitor-agent --config=/etc/monitor.ini
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
EOF